# crontab/cron.sh
#!/usr/bin/python
10 10 * * 1-5 python /Users/xiaxiaodong/python/python3/crontab/start_spider.py
